# launchd/com.sotaurabe.arxiv-slack-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sotaurabe.arxiv-slack-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/user/arxiv-research-digest/scripts/run_slack_agent.sh</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>300</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/Users/user/.local/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/user/arxiv-research-digest/launchd_agent_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/arxiv-research-digest/launchd_agent_stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
