# launchd/com.sotaurabe.arxiv-weekly-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sotaurabe.arxiv-weekly-digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/user/arxiv-research-digest/scripts/weekly_slack_digest.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/Users/user/.local/bin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/user/arxiv-research-digest/launchd_weekly_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/arxiv-research-digest/launchd_weekly_stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
